# scripts/script/10-recover.rules
SUBSYSTEM=="block", SUBSYSTEMS=="usb", ACTION=="add", KERNEL=="sd??", RUN+="/usr/bin/mkarchive ADD %k %s{size}"
